# systemd/litclock-update.service
[Unit]
Description=LitClock Weekly Auto-Update
Documentation=https://github.com/kapoorankush/litclock/issues/209
# Don't start mid-firstboot — the provisioning flow owns the clock at that
# point and a concurrent git reset would race setup state. The real
# "firstboot finished" signal is /etc/litclock/.setup-complete (written
# by first-boot.sh and read by boot-splash.sh).
After=network-online.target litclock-firstboot.service
Wants=network-online.target
ConditionPathExists=/etc/litclock/.setup-complete

[Service]
Type=oneshot
ExecStart=/home/pi/litclock/scripts/update.sh
WorkingDirectory=/home/pi/litclock
User=pi
# Bound the worst case for a real update on a Pi Zero 2W. M5 plan-eng-review
# (#245) caught the prior 120s ceiling: a real pip-install + image-sync + smoke
# can exceed 2 min under CPU pressure, so the timer was killing real updates
# mid-Phase-4. 600s lets a worst-case run finish; the script's own internal
# `timeout 60` on smoke + `timeout 30` on git-fetch is the second layer.
# F1 in PLAN-LitClock-Control-PWA.md "Codex amendments" table.
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
