# systemd/litclock-reresolve-location.service
[Unit]
Description=LitClock on-boot location re-resolve (#337 A2/A8)
Documentation=https://github.com/kapoorankush/litclock/issues/337
# #337 A2: gate on the handoff-complete marker. On a fresh-flash boot, setup_server
# is still doing its own post-WiFi IP-geo resolve via _resolve_location_from_ip;
# firing this oneshot in parallel would do a redundant ip-api.com call AND fight
# the same env.sh flock. The marker is only written once setup + handoff fully
# complete — by then setup_server is gone and this oneshot is the right caller.
# On subsequent boots the marker persists, so this oneshot fires normally.
ConditionPathExists=/etc/litclock/.handoff-complete
# #337 A8 (deliberate): NO `Before=litclock.service`. Best-effort boot. If IP-geo
# is slow (DNS race post-WiFi, ip-api outage), the first quote tick still fires
# on schedule using last-boot's data; the next tick has the fresh result.
# Blocking the first quote on a 1–33s IP-geo would be a real boot-delay regression.
After=NetworkManager-wait-online.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Python entry point in src/location_resolver.py. Reads env.sh, gates on
# WEATHER_LOCATION_MODE=auto (CRITICAL silent-corruption guard per A15 — a user
# in Specific mode has chosen their location intentionally, this MUST NOT
# overwrite it), and runs the full-retry IP-geo path (1+3+9s backoff, ~33s
# worst case). Writes lat/lon/name/units/tz/MODE/IP_COUNTRY atomically via
# config.atomic_update — same flock as every other writer.
WorkingDirectory=/home/pi/litclock
EnvironmentFile=/home/pi/litclock/env.sh
# LITCLOCK_ENV_FILE tells location_resolver.main() which file to read/write.
# Without this the oneshot would fall through to setup_server.ENV_FILE which
# is None outside the first-boot context.
Environment=LITCLOCK_ENV_FILE=/home/pi/litclock/env.sh
# Script path (NOT `python -m location_resolver`). `python -m` searches cwd +
# venv site-packages; `src/` is on neither path on the Pi (PYTHONPATH is not
# wired up for this unit, and `src/` is registered only for pytest via
# pyproject.toml's `pythonpath`). The script-path form makes Python add the
# script's dir to sys.path, matching the existing pattern used by
# litclock-control.service for control_server/app.py.
ExecStart=/home/pi/litclock/venv/bin/python3 /home/pi/litclock/src/location_resolver.py
User=pi
StandardOutput=journal
StandardError=journal
# Tight timeout — the resolver's own retry budget caps at ~33s, plus a small
# margin for the env.sh flock contention window (#274). Anything past 60s
# means something is wrong (DNS wedged?) and we want systemd to surface it.
TimeoutStartSec=60

[Install]
# Run once at boot — the systemd convention for a oneshot that should fire
# automatically without an explicit trigger. NetworkManager-wait-online (in
# After=) gates the actual execution until the network is up.
WantedBy=multi-user.target
